# deploy/linux/clouda-cleanup.timer
[Unit]
Description=Daily Clouda temporary cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
