# deploy/linux/clouda-api.service
[Unit]
Description=Clouda internal worker API
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=clouda
Group=clouda
WorkingDirectory=/opt/clouda
EnvironmentFile=-/opt/clouda/.env
ExecStart=/opt/clouda/deploy/linux/run_api.sh
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/clouda/data /opt/clouda/conversions /opt/clouda/logs

[Install]
WantedBy=multi-user.target
